# iso/config/includes.chroot/etc/systemd/system/foundry.service
[Unit]
Description=Foundry IPTV
After=docker.service foundry-firstboot.service
Requires=docker.service
Wants=foundry-firstboot.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/foundry
ExecStartPre=/bin/bash -c 'test -f /opt/foundry/.firstboot-done'
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose up -d --remove-orphans
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
